# airootfs/etc/systemd/system/sshd.service
# This unit file intentionally masks the OpenSSH daemon in the live ISO.
# Systemd will treat this as a non-runnable placeholder (similar to a /dev/null mask)
# because the unit always fails its Condition and has no meaningful ExecStart.
#
# Path: ArchISO/airootfs/etc/systemd/system/sshd.service

[Unit]
Description=OpenSSH Daemon (masked in ArchRiot live ISO)
Documentation=man:sshd(8) man:sshd_config(5)
# Never start in the live environment
Condition=false

[Service]
Type=oneshot
# Neutralize any original ExecStart from the packaged unit
ExecStart=
ExecStart=/bin/true
RemainAfterExit=yes

[Install]
# No targets should ever want this in the live ISO
WantedBy=
