# airootfs/etc/systemd/system/sshd.socket.d/disable.conf
[Unit]
# Disable sshd socket activation in the live ISO environment to prevent sshd from starting.
Condition=false
